# rtaudio/target/rtaudio-wifi.service
[Unit]
Description=Create AP Service
After=network.target

[Service]
Type=simple
ExecStart=/home/rtaudio/rtaudio-wifi.sh
KillSignal=SIGINT
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
