# root/var/lib/pacman/local/hostapd-2.5-1/files
%FILES%
etc/
etc/hostapd/
usr/
usr/bin/
usr/bin/hostapd
usr/bin/hostapd_cli
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/hostapd.service
usr/share/
usr/share/doc/
usr/share/doc/hostapd/
usr/share/doc/hostapd/hlr_auc_gw.milenage_db
usr/share/doc/hostapd/hostapd.accept
usr/share/doc/hostapd/hostapd.conf
usr/share/doc/hostapd/hostapd.deny
usr/share/doc/hostapd/hostapd.eap_user
usr/share/doc/hostapd/hostapd.radius_clients
usr/share/doc/hostapd/hostapd.sim_db
usr/share/doc/hostapd/hostapd.vlan
usr/share/doc/hostapd/hostapd.wpa_psk
usr/share/doc/hostapd/wired.conf
usr/share/licenses/
usr/share/licenses/hostapd/
usr/share/licenses/hostapd/COPYING
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/hostapd_cli.1.gz
usr/share/man/man8/
usr/share/man/man8/hostapd.8.gz

